# bybit-bot.service
# Bybit Demo Trading Bot - Systemd Service
# 
# Installation:
#   1. Copy: sudo cp bybit-bot.service /etc/systemd/system/elcaro-bot.service
#   2. Reload: sudo systemctl daemon-reload
#   3. Enable: sudo systemctl enable elcaro-bot
#   4. Start: sudo systemctl start elcaro-bot
#
# Commands:
#   sudo systemctl status elcaro-bot
#   sudo systemctl start elcaro-bot
#   sudo systemctl stop elcaro-bot
#   sudo systemctl restart elcaro-bot
#   journalctl -u elcaro-bot -f

[Unit]
Description=Elcaro Bybit Trading Bot v2
After=network.target
# Don't conflict with other services
Conflicts=

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/project/elcarobybitbotv2
Environment=PATH=/home/ubuntu/project/elcarobybitbotv2/venv/bin:/usr/bin:/bin
ExecStart=/home/ubuntu/project/elcarobybitbotv2/venv/bin/python bot.py
Restart=on-failure
RestartSec=10

# Don't restart too fast (prevents infinite restart loops)
StartLimitIntervalSec=300
StartLimitBurst=5

# Logging to file
StandardOutput=append:/home/ubuntu/project/elcarobybitbotv2/bot.log
StandardError=append:/home/ubuntu/project/elcarobybitbotv2/bot.log

# Security - isolated from other services
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
